# systemd/liveuamap-scrapes@.service
[Unit]
Description=Liveuamap scrapes for %i-archive
OnFailure=status-email-archive@%n.service

[Service]
Type=oneshot
ExecStart=/bin/sh /home/%i/%i-archive/bin/liveuamap-daily.sh
User=%i
WorkingDirectory=/home/%i/%i-archive

[Install]
RequiredBy=default.target
