# systemd/youtube-videos-scrapes@.timer
[Unit]
Description=Run Youtube videos scrapes daily

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
