# systemd/youtube-historic-scrapes@.timer
[Unit]
Description=Run Youtube historic scrapes daily

[Timer]
OnCalendar=*-*-* 01:00:00
Persistent=true

[Install]
WantedBy=timers.target
